# opt/cornell-tilde/systemd/cornell-tilde-directory.path
[Unit]
Description=Watch Cornell Tilde database for relevant directory modifications

[Path]
PathChanged=/opt/cornell-tilde/var/cornell_tilde.sqlite3
Unit=cornell-tilde-directory.service

[Install]
WantedBy=multi-user.target